# debian/postinst
#!/bin/sh
# postinst script
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package

case "$1" in
    configure)
    ;;
    abort-upgrade|abort-remove|abort-deconfigure)
    ;;
    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

# Enable default templates
if [ "$1" = "configure" ] || [ "$1" = "abort-upgrade" ] || [ "$1" = "abort-deconfigure" ] || [ "$1" = "abort-remove" ] ; then
	systemctl unmask 'docker-compose@etc-docker\x2dcompose' >/dev/null || true
	systemctl enable 'docker-compose@etc-docker\x2dcompose' >/dev/null || true
	systemctl unmask 'docker-compose-refresh@etc-docker\x2dcompose' >/dev/null || true
	systemctl enable 'docker-compose-refresh@etc-docker\x2dcompose' >/dev/null || true
	if [ -d /run/systemd/system ]; then
		systemctl --system daemon-reload >/dev/null || true
		if [ -n "$2" ]; then
			_dh_action=restart
		else
			_dh_action=start
		fi
		systemctl $_dh_action 'docker-compose@etc-docker\x2dcompose' >/dev/null || true
		systemctl $_dh_action 'docker-compose-refresh@etc-docker\x2dcompose' >/dev/null || true
	fi
fi

exit 0

